# supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm -dmemory_limit=1024M --allow-to-run-as-root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
